# yaat-sidecar.service
[Unit]
Description=YAAT Sidecar - Backend Monitoring Agent
Documentation=https://docs.yaat.io
After=network.target

[Service]
Type=simple
User=yaat
Group=yaat
WorkingDirectory=/var/lib/yaat
ExecStart=/usr/local/bin/yaat-sidecar --config /etc/yaat/yaat.yaml --log-file /var/log/yaat/sidecar.log
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/yaat /var/lib/yaat
CapabilityBoundingSet=

# Resource limits
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
